# .config/systemd/user/vdirsyncer.timer
[Unit]
Description=Calendar synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=2m
Unit=vdirsyncer.service

[Install]
WantedBy=timers.target
